# systemd/sonic-troubleshooting-demo-service.service
[Unit]
Description=SONiC Troubleshooting Demo Service
After=database.service
Requires=database.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sonic-troubleshooting-demo-service remediate --port Ethernet0 --target-status down --write-mode cli --restore --hold-seconds 5
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
